# configs/supervisor.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=10
stdout_logfile=/workspace/logs/cron/stdout.log
stderr_logfile=/workspace/logs/cron/stderr.log
